# conf/etc/cron.hourly/logrotate
#!/bin/sh

/usr/sbin/logrotate /etc/logrotate.hourly.conf